# crontab
SHELL=/bin/sh
*/10 * * * * /script.sh >> /var/log/script.log
